# etc/services.d/molly-brown/run
#!/usr/bin/with-contenv sh

exec s6-setuidgid molly-brown molly-brown -c /etc/molly-brown/molly.conf
